# copy-files/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:tidal-scheduler]
command=python /app/scheduler.py
directory=/app
autostart=true
autorestart=true
environment=TZ="%(ENV_TZ)s",CRON_SCHEDULE="%(ENV_CRON_SCHEDULE)s"
priority=100
startsecs=5
startretries=3
stderr_logfile=/var/log/supervisor/tidal-scheduler-stderr.log
stdout_logfile=/var/log/supervisor/tidal-scheduler-stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
